# server/supervisor/development.ini
[program:mybic-development]
directory = /home/thedevel/sites/mybic-dev-env/mybic
command = /home/thedevel/sites/mybic-dev-env/bin/uwsgi --ini /home/thedevel/sites/mybic-dev-env/mybic/server/uwsgi/development.ini
autostart = true
autorestart = true
