# server/supervisor/development.ini
[program:mybic-development]
directory = /home/devuser/webapps/mybic-dev-env/mybic-dev
command = /home/devuser/webapps/mybic-dev-env/bin/uwsgi --ini /home/devuser/webapps/mybic-dev-env/mybic-dev/server/uwsgi/debug.ini --uid devuser --gid devuser --py-auto-reload 3
numprocs=1
directory=/tmp
umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=devuser
redirect_stderr=true
stdout_logfile=/var/log/mybic-dev.stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stderr_logfile=/var/log/mybic-dev.stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
serverurl=AUTO

[program:solr-development]
user = solr
directory = /opt/solr
command = /usr/java/java7/bin/java -Xmx2500m -jar start.jar
