# server/supervisor/production.ini
[program:mybic-production]
directory = /home/devuser/webapps/mybic-env/mybic
command = /home/devuser/webapps/mybic-env/bin/uwsgi --ini /home/devuser/webapps/mybic-env/mybic/server/uwsgi/production.ini --uid mybic --gid dbhi_bic
numprocs=1
directory=/tmp
umask=022
priority=999
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=mybic
redirect_stderr=true
stdout_logfile=/var/log/mybic-prod.stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
stderr_logfile=/var/log/mybic-prod.stderr.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_capture_maxbytes=1MB
serverurl=AUTO
