# server/supervisor/production.ini
[program:mybic-production]
directory = /home/thedevel/sites/mybic-env/mybic
command = /home/thedevel/sites/mybic-env/bin/uwsgi --ini /home/thedevel/sites/mybic-env/mybic/server/uwsgi/production.ini
autostart = true
autorestart = true
